# crontab
*/15 * * * * root php -q /app/cron_temp_log.php
*/5 * * * * root php -q /app/cron_schedule.php
*/1 * * * * root echo hypriot/rpi-cron-example >>/var/log/cron.log 2>&1
